# assets/config/systemd/nekotun.service
[Unit]
Description=Neko Tunneling Backend API
Documentation=https://github.com/xenvoid404/autoscript-tunneling
After=network.target nss-lookup.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/node /opt/nekotun/dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
